# gluon-usb-media/files/etc/hotplug.d/block/a-20-usb-mount-media
USBMEDIAPATH=$(uci get 'usb-media.settings.path')
case "$ACTION" in
    add)
        if [ "1" == $(uci get 'usb-media.settings.share_devices') ]; then
          #if $DEVICENAME begins with "s"
          case "$DEVICENAME" in
               s*)
                   if [ "0" == $(grep $DEVICENAME /proc/mounts|echo $?) ]; then
                     [ -d /mnt/$DEVICENAME ] || mkdir /mnt/$DEVICENAME
                     logger /mnt/$DEVICENAME was created
                     mount -o ro /dev/$DEVICENAME /mnt/$DEVICENAME || rmdir /mnt/$DEVICENAME
                     logger the media was mounted
                     [ -d /lib/gluon/status-page/www/$USBMEDIAPATH ] || mv /lib/gluon/status-page/www/media /lib/gluon/status-page/www/$USBMEDIAPATH
                   fi
                   ;;
          esac
        fi
        ;;
    remove)
        logger "Remove $DEVICENAME"
        [ -d /mnt/$DEVICENAME ] && rmdir /mnt/$DEVICENAME
        [ -L /lib/gluon/status-page/www/media ] || mv /lib/gluon/status-page/www/$USBMEDIAPATH /lib/gluon/status-page/www/media
        ;;
esac
